# deploy/systemd/ombudsman-backend.service
[Unit]
Description=Ombudsman Validation Studio - Backend API
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/data/ombudsman/ombudsman-validation-studio/backend
EnvironmentFile=/data/ombudsman/ombudsman.env

# Set Python path
Environment="PYTHONPATH=/data/ombudsman/ombudsman-validation-studio/backend:/data/ombudsman/ombudsman_core/src"
Environment="OMBUDSMAN_DATA_DIR=/data/ombudsman/data"
Environment="OMBUDSMAN_CORE_DIR=/data/ombudsman/ombudsman_core/src/ombudsman/config"
Environment="OMBUDSMAN_LOG_DIR=/data/ombudsman/logs"

ExecStart=/data/ombudsman/ombudsman-validation-studio/backend/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000 --log-level info

Restart=always
RestartSec=5

# Logging
StandardOutput=append:/data/ombudsman/logs/backend.log
StandardError=append:/data/ombudsman/logs/backend.log

[Install]
WantedBy=multi-user.target
