# deploy/systemd/ombudsman-backend.service
[Unit]
Description=Ombudsman Validation Studio Backend (FastAPI)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=ombudsman
Group=ombudsman
WorkingDirectory=/opt/ombudsman/backend

# Environment file
EnvironmentFile=/opt/ombudsman/.env

# Python path setup
Environment="PYTHONPATH=/opt/ombudsman/backend:/opt/ombudsman/core/src"

# Data storage paths
Environment="MAPPING_INTELLIGENCE_DIR=/var/lib/ombudsman/data/mapping_intelligence"
Environment="QUERY_HISTORY_DIR=/var/lib/ombudsman/data/query_history"
Environment="PIPELINE_RUNS_DIR=/var/lib/ombudsman/data/pipeline_runs"
Environment="CONFIG_BACKUPS_DIR=/var/lib/ombudsman/data/config_backups"
Environment="AUTH_DATA_DIR=/var/lib/ombudsman/data/auth"
Environment="AUDIT_LOG_DIR=/var/lib/ombudsman/data/audit_logs"
Environment="BATCH_JOBS_DIR=/var/lib/ombudsman/data/batch_jobs"
Environment="NOTIFICATION_RULES_FILE=/var/lib/ombudsman/data/notification_rules.json"

# Ollama URL for native Linux (localhost instead of host.docker.internal)
Environment="OLLAMA_BASE_URL=http://localhost:11434"

# Start command
ExecStart=/opt/ombudsman/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 4

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/ombudsman/backend.log
StandardError=append:/var/log/ombudsman/backend-error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ombudsman /var/log/ombudsman /opt/ombudsman/backend/data

[Install]
WantedBy=multi-user.target
