# deploy/systemd/ombudsman-backend.service
[Unit]
Description=Ombudsman Validation Studio Backend (FastAPI)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=ombudsman
Group=ombudsman
WorkingDirectory=/opt/ombudsman/backend

# Environment file
EnvironmentFile=/opt/ombudsman/.env

# Python path setup
Environment="PYTHONPATH=/opt/ombudsman/backend:/opt/ombudsman/core/src"

# Centralized path configuration - ALL paths derive from these two variables
# Set ENVIRONMENT=production to use production paths
Environment="ENVIRONMENT=production"
Environment="OMBUDSMAN_DATA_DIR=/var/lib/ombudsman/data"
Environment="OMBUDSMAN_CORE_DIR=/opt/ombudsman/core/src/ombudsman/config"
Environment="OMBUDSMAN_LOG_DIR=/var/log/ombudsman"

# Ollama URL for native Linux (localhost instead of host.docker.internal)
Environment="OLLAMA_BASE_URL=http://localhost:11434"

# Start command
ExecStart=/opt/ombudsman/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 4

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/ombudsman/backend.log
StandardError=append:/var/log/ombudsman/backend-error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ombudsman /var/log/ombudsman /opt/ombudsman/backend/data

[Install]
WantedBy=multi-user.target
